# services/llmspaghetti-terminal.service
[Unit]
Description=LLMSpaghetti Web Terminal (ttyd)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as llmspaghetti user — restricts access to that account by default.
# Power users can sudo from this shell just like any other terminal.
User=llmspaghetti
Group=llmspaghetti

# ttyd flags (ttyd 1.7.x):
#   --port 7681      listen port (proxied by Caddy at /terminal/)
#   --interface lo   localhost only — Caddy handles auth + TLS
#   --writable       allow clients to write (required for input)
#   -t fontSize=15   default font size
#   -t theme=...     dark theme matching LLMSpaghetti UI
#   /bin/bash        the shell to launch
# NOTE: no --reconnect flag — ttyd 1.7.x doesn't have it; it breaks arg parsing.
ExecStart=/usr/bin/ttyd \
  --port 7681 \
  --interface lo \
  --writable \
  -t fontSize=15 \
  -t fontFamily="'JetBrains Mono','Fira Code','Cascadia Code',monospace" \
  -t theme='{"background":"#0d1117","foreground":"#e6edf3","cursor":"#58a6ff","black":"#484f58","red":"#f85149","green":"#3fb950","yellow":"#d29922","blue":"#58a6ff","magenta":"#bc8cff","cyan":"#39c5cf","white":"#b1bac4","brightBlack":"#6e7681","brightRed":"#ff7b72","brightGreen":"#56d364","brightYellow":"#e3b341","brightBlue":"#79c0ff","brightMagenta":"#d2a8ff","brightCyan":"#56d4dd","brightWhite":"#f0f6fc"}' \
  /bin/bash --login

Restart=always
RestartSec=3
StandardOutput=append:/opt/llmspaghetti/logs/ttyd.log
StandardError=append:/opt/llmspaghetti/logs/ttyd.log

[Install]
WantedBy=multi-user.target
