# services/llmspaghetti-firstboot.service
[Unit]
Description=LLMSpaghetti First-Boot Setup Wizard
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/opt/llmspaghetti/.firstboot-complete

[Service]
Type=simple
User=root
WorkingDirectory=/opt/llmspaghetti/firstboot
ExecStart=/opt/llmspaghetti/.venv/bin/python /opt/llmspaghetti/firstboot/main.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/llmspaghetti/logs/firstboot.log
StandardError=append:/opt/llmspaghetti/logs/firstboot.log

[Install]
WantedBy=multi-user.target
